# config/supervisor/supervisord.conf
[supervisord]
nodaemon=true
loglevel=debug

[program:sshd]
command=/usr/sbin/sshd
autorestart=true

[program:php8-fpm]
command=/usr/sbin/php-fpm8.0
autorestart=true

[program:cron]
command=cron -f
